# templates/default/init.erb
#!/bin/bash 
### BEGIN INIT INFO
# Provides:          kibana
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Make sense of a mountain of logs.
### END INIT INFO

#
# chkconfig: 2345 20 80
#
export PID_DIR="<%= node['kibana']['pid_dir'] %>"
KIBANA_HOME="<%= @basedir %>"
KIBANA_BIN="kibana-daemon.rb" #short name
KIBANA_SRV="$KIBANA_HOME/$KIBANA_BIN" #script to issue start/stop/etc commands
PATH="<%= node['rbenv']['root_path'] %>/shims/:$PATH"
KIBANA_USER="<%= node['kibana']['user'] %>"
#LOGDIR="<%= node['kibana']['log_dir'] %>"

function missing_bin() { 
    echo "Failed to find kibana binary file: [$KSRV]"; 
    echo "Check /etc/init.d/kibana file for correct settings."
    RETVAL=1; 
    exit 1;
}

which ruby > /dev/null
if [ $? -ne 0 ]; then 
    echo "Ruby cannot be found. Please install ruby or put in PATH."
    exit 1;
fi

test -f $KIBANA_SRV || missing_bin

. /etc/rc.d/init.d/functions

RETVAL=0

case "$1" in
start)
  echo -n "Starting Kibana: "

  #check to see if we're already running
  pgrep -f ${KIBANA_BIN} > /dev/null
  RUNNING=$?
  if [ $RUNNING -eq 0 ]; then             
      echo "[FAILED]"
      echo        
      echo "Reason: kibana is already running."
      RETVAL=1
      exit 1;
  fi
  su - $KIBANA_USER -c "ruby $KIBANA_SRV $1"
  RETVAL=$?
;;

stop)
echo -n "Shutting down Kibana: "
  su - $KIBANA_USER -c "ruby $KIBANA_SRV $1"
  RETVAL=$?
;;

restart)
  su - $KIBANA_USER -c "ruby $KIBANA_SRV $1"
  RETVAL=$?
;;

status)
  su - $KIBANA_USER -c "ruby $KIBANA_SRV $1"
;;

# zap not working yet
#force-reload)
#  su - $KIBANA_USER -c "ruby $KIBANA_SRV zap"
#  su - $KIBANA_USER -c "ruby $KIBANA_SRV start"
#  RETVAL=$?
#;;

*)
  echo "Usage: $0 {start|stop|restart|status}"
  exit 1
esac

exit $RETVAL